# hotspot.service
[Unit]
Description=Hotspot Pix Flask App
After=network.target

[Service]
WorkingDirectory=/home/pi/hotspot-pix
ExecStart=/usr/bin/make run
ExecStop=/usr/bin/make stop
Restart=always
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

# --- systemd service: hotspot.service ---

[Unit]
Description=Hotspot Flask Server
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/udoo/hotspot-pix/app.py
WorkingDirectory=/home/udoo/hotspot-pix
Restart=always
Environment=FLASK_ENV=production
EnvironmentFile=/home/udoo/hotspot-pix/.env

[Install]
WantedBy=multi-user.target